# superset/supervisord.conf.d.conf

[supervisord]
nodaemon=true
user=superset

[program:sshd]
command=/usr/sbin/sshd -f /home/superset/custom_ssh/sshd_config
autostart=true
redirect_stderr=True